# crontab
0 0 * * * cd /var/www/fontbakery && python scripts/statupdate.py
